# bmc-network-online.target.in
[Unit]
Description=Bmc network online target
